# piper-yang/rootfs/etc/s6-overlay/s6-rc.d/discovery/up
00000000: 6f6e 6573 686f 740a                      oneshot.
