# piper-yang/rootfs/etc/s6-overlay/s6-rc.d/user/run
#!/command/with-contenv bashio
# vim: ft=bash
# shellcheck shell=bash
exec s6-svscan /var/run/s6/services
